# lib/systemd/system/monitor@.service
[Unit]
Description=Monitor for %i

[Service]
ExecStart = /var/lib/ybs32/bin/monitor.sh %i
SyslogIdentifier=monitor
