# lib/systemd/system/unico.service
[Unit]
Description=ASGI server
OnFailure=monitor@%n.service

[Service]
ExecStart=/var/lib/ybs32/bin/unico.sh
SyslogIdentifier=unico
Restart = on-failure

[Install]
WantedBy=multi-user.target
