# etc/systemd/system/monitor@.service
[Unit]
Description=Monitor for %i

[Service]
ExecStart = /opt/ybs32/bin/monitor.sh %i
SyslogIdentifier=monitor
